# config/examples/com.openclaw.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.openclaw.watchdog</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>__REPO_ROOT__/scripts/start_watchdog.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>__REPO_ROOT__</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/openclaw-watchdog.out.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/openclaw-watchdog.err.log</string>
  </dict>
</plist>
